# install/udev/99-retrosync-fxpak.rules
# RetroSync — poke the daemon for an immediate sync when the FXPak Pro
# (or any usb2snes-compatible cart sharing the same USB IDs) is connected.
#
# Without this, cart-on → first sync latency is bounded by the
# orchestrator's `unhealthy_recheck_sec` (default 2s). With it, the
# daemon does an immediate pass on USB attach (sub-second latency).
#
# Vendor:product 1209:5a22 is "Generic ikari_01 sd2snes" — pid.codes
# allocation used by the FXPak Pro / sd2snes firmware. Confirm yours
# with `lsusb` if the rule doesn't fire (some clones may use different
# IDs).
#
# The rule signals the running retrosync daemon via SIGUSR1, which
# pokes every orchestrator into running an immediate pass.

ACTION=="add", SUBSYSTEM=="tty", \
  ATTRS{idVendor}=="1209", ATTRS{idProduct}=="5a22", \
  RUN+="/bin/systemctl kill --signal=SIGUSR1 retrosync.service"

# Fallback for older udev that doesn't recurse to the parent USB
# device's ATTRS at the tty level.
ACTION=="add", SUBSYSTEM=="usb", \
  ATTR{idVendor}=="1209", ATTR{idProduct}=="5a22", \
  RUN+="/bin/systemctl kill --signal=SIGUSR1 retrosync.service"
