# install/systemd/retrosync.service
[Unit]
Description=RetroSync save backup daemon
Documentation=https://github.com/indiefan/RetroSync
After=network-online.target sni.service
Wants=network-online.target sni.service

[Service]
Type=simple
User=retrosync
Group=retrosync
Environment=RETROSYNC_CONFIG=/etc/retrosync/config.yaml
# rclone's default config search path is ~/.config/rclone/rclone.conf, which
# is hidden by ProtectHome=true. Even adding /home/retrosync/.config/rclone
# to ReadWritePaths doesn't make the parent directory traversable inside
# the namespace systemd builds. We keep the rclone config under /var/lib
# (already in the writable space) instead. Both the daemon and the CLI
# read from this exact path; setup.sh writes the OAuth credentials there.
Environment=RCLONE_CONFIG=/var/lib/retrosync/rclone.conf
ExecStart=/opt/retrosync/.venv/bin/retrosyncd --config ${RETROSYNC_CONFIG}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/retrosync /home/retrosync/.config/rclone

[Install]
WantedBy=multi-user.target
